# conf/systemd/gogdb-fetch-backup.timer
[Unit]
Description=GOG DB Fetch Backup Timer

[Timer]
OnCalendar=00:10:15
AccuracySec=5min

[Install]
WantedBy=timers.target
