# conf/systemd/gogdb-backup.timer
[Unit]
Description=GOG DB Backup Timer

[Timer]
OnCalendar=00:01:15
AccuracySec=1min

[Install]
WantedBy=timers.target
